# ansible/roles/prometheus/templates/prometheus.service.j2
[Unit]
Description=Prometheus Monitoring System (User Service)
Documentation=https://prometheus.io/docs/introduction/overview/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ monitoring_service_user }}
Group={{ monitoring_group }}

# Рабочая директория
WorkingDirectory={{ monitoring_base_dir }}

# Исполняемый файл и аргументы
ExecStart={{ monitoring_dirs.bin }}/prometheus \
  --config.file={{ monitoring_dirs.config }}/prometheus.yml \
  --storage.tsdb.path={{ monitoring_dirs.data }}/prometheus \
  --web.console.templates={{ monitoring_dirs.config }}/consoles \
  --web.console.libraries={{ monitoring_dirs.config }}/console_libraries \
  --web.config.file={{ monitoring_dirs.config }}/web-config.yml \
  --web.external-url=https://%H:{{ prometheus_port }}/ \
  --web.listen-address=0.0.0.0:{{ prometheus_port }}

# Graceful shutdown
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=30s
Restart=on-failure
RestartSec=10s

# Переменные окружения
Environment="PATH={{ monitoring_dirs.bin }}:/usr/local/bin:/usr/bin:/bin"

# Ограничения ресурсов
LimitNOFILE=65536
LimitNPROC=512

# Безопасность
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths={{ monitoring_dirs.data }}
ReadWritePaths={{ monitoring_dirs.logs }}
ReadOnlyPaths={{ monitoring_dirs.bin }}
ReadOnlyPaths={{ monitoring_dirs.config }}
ReadOnlyPaths=/dev/shm/monitoring_secrets

# Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=prometheus

[Install]
WantedBy=default.target

