# ansible/roles/rlm_standard_setup/templates/prometheus.service.j2
[Unit]
Description=Prometheus Monitoring System
Documentation=https://prometheus.io/docs/introduction/overview/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=prometheus
Group=prometheus

# Рабочая директория
WorkingDirectory=/var/lib/prometheus

# Исполняемый файл и аргументы
ExecStart=/usr/bin/prometheus \
  --config.file=/etc/prometheus/prometheus.yml \
  --web.config.file=/etc/prometheus/web-config.yml \
  --storage.tsdb.path=/var/lib/prometheus/data \
  --web.console.templates=/etc/prometheus/consoles \
  --web.console.libraries=/etc/prometheus/console_libraries \
  --web.listen-address=0.0.0.0:{{ prometheus_port }}

# Graceful shutdown
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=30s
Restart=on-failure
RestartSec=10s

# Переменные окружения
EnvironmentFile=-/etc/prometheus/prometheus.env

# Ограничения ресурсов
LimitNOFILE=10000
LimitNPROC=512

# Безопасность
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/prometheus
ReadOnlyPaths=/etc/prometheus

# Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=prometheus

[Install]
WantedBy=multi-user.target
