# ansible/roles/harvest/templates/harvest.service.j2
[Unit]
Description=NetApp Harvest Metrics Collector (User Service)
Documentation=https://netapp.github.io/harvest/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Рабочая директория
WorkingDirectory={{ monitoring_base_dir }}

# Исполняемый файл и аргументы
ExecStart={{ monitoring_dirs.bin }}/harvest \
  --config {{ monitoring_dirs.config }}/harvest.yml \
  --loglevel info \
  --promPort {{ harvest_netapp_port }} \
  --restPort {{ harvest_unix_port }}

# Graceful shutdown
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=30s
Restart=on-failure
RestartSec=10s

# Переменные окружения
Environment="PATH={{ monitoring_dirs.bin }}:/usr/local/bin:/usr/bin:/bin"
Environment="HARVEST_HOME={{ monitoring_base_dir }}"

# Ограничения ресурсов
LimitNOFILE=10000
LimitNPROC=512

# Безопасность
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths={{ monitoring_dirs.data }}/harvest
ReadWritePaths={{ monitoring_dirs.logs }}/harvest
ReadOnlyPaths={{ monitoring_dirs.bin }}
ReadOnlyPaths={{ monitoring_dirs.config }}
ReadOnlyPaths=/dev/shm/monitoring_secrets

# Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=harvest

[Install]
WantedBy=default.target



