# ansible/roles/harvest/templates/harvest.service.j2
[Unit]
Description=NetApp Harvest Metrics Collector (User Service)
Documentation=https://netapp.github.io/harvest/
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User={{ monitoring_service_user }}
Group={{ monitoring_group }}

# Рабочая директория
WorkingDirectory={{ monitoring_base_dir }}

# Запуск всех pollers
ExecStart={{ monitoring_dirs.bin }}/harvest start --config {{ monitoring_dirs.config }}/harvest.yml

# Остановка всех pollers
ExecStop={{ monitoring_dirs.bin }}/harvest stop --config {{ monitoring_dirs.config }}/harvest.yml

# Перезапуск
ExecReload={{ monitoring_dirs.bin }}/harvest restart --config {{ monitoring_dirs.config }}/harvest.yml

# Таймауты
TimeoutStartSec=60s
TimeoutStopSec=60s
Restart=on-failure
RestartSec=15s

# Переменные окружения
Environment="PATH={{ monitoring_dirs.bin }}:/usr/local/bin:/usr/bin:/bin"
Environment="HARVEST_CONF={{ monitoring_dirs.config }}/harvest.yml"

# Ограничения ресурсов
LimitNOFILE=65536
LimitNPROC=1024

# Безопасность
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths={{ monitoring_dirs.logs }}/harvest
ReadOnlyPaths={{ monitoring_dirs.bin }}
ReadOnlyPaths={{ monitoring_dirs.config }}
ReadOnlyPaths=/dev/shm/monitoring_secrets

# Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=harvest

[Install]
WantedBy=default.target



